# freshrss/rootfs/etc/s6-overlay/s6-rc.d/init-nginx/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: FreshRSS
# Configures NGINX for use with freshrss
# ==============================================================================

# Generate Ingress PHP-FPM configuration
INTERFACE=$(bashio::addon.ip_address)
INGRESS_PATH=$(bashio::addon.ingress_entry)
INGRESS_PORT=$(bashio::addon.ingress_port)

bashio::log.debug "Using interface: ${INTERFACE}"
bashio::log.debug "Using ingress path: ${INGRESS_PATH}"
bashio::log.debug "Using ingress port: ${INGRESS_PORT}"

bashio::var.json \
    interface "${INTERFACE}" \
    ingress_path "${INGRESS_PATH}" \
    ingress_port "${INGRESS_PORT}" \
    | tempio \
        -template /etc/nginx/templates/ingress.gtpl \
        -out /etc/nginx/servers/ingress.conf

# Generate direct access configuration, if enabled.
if bashio::var.has_value "$(bashio::addon.port 80)"; then
    bashio::config.require.ssl
    bashio::var.json \
        certfile "$(bashio::config 'certfile')" \
        keyfile "$(bashio::config 'keyfile')" \
        ssl "^$(bashio::config 'ssl')" \
        | tempio \
            -template /etc/nginx/templates/direct.gtpl \
            -out /etc/nginx/servers/direct.conf
fi
