# freshrss/rootfs/etc/s6-overlay/s6-rc.d/freshrss/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: FreshRSS
# Runs Refresh script
# ==============================================================================

bashio::log.info "Starting Refresh..."

exec /usr/bin/refresh.sh