# freshrss/rootfs/etc/s6-overlay/s6-rc.d/php-fpm/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Add-on: FreshRSS
# Runs the PHP-FPM daemon
# ==============================================================================
bashio::log.info "Starting PHP-FPM..."

exec php-fpm81 -R --nodaemonize